# udev/99-usb-arduino-uno-serial.rules
ACTION=="add", SUBSYSTEM=="tty", ATTRS{serial}=="1234567890ABCDEFGH0002", RUN+="/bin/sh -c 'stty -F $devnode 2000000 raw -icrnl"
